# systemd/wallpaper.timer
[Unit]
Description=Run wallpaper shuffler every X minutes

[Timer]
OnBootSec=1
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
