# wallpaper_changer/wallpaper.timer
[Unit]
Description=Run wallpaper shuffler every X minutes

[Timer]
OnBootSec=30min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
